# deployment/crimecity3k.service
[Unit]
Description=CrimeCity3K - Swedish Police Events Map
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=10s
TimeoutStartSec=60s
TimeoutStopSec=30s

# Remove existing container if present (handles upgrades)
ExecStartPre=-/usr/bin/podman stop crimecity3k
ExecStartPre=-/usr/bin/podman rm crimecity3k

# Start container
# - Uses :production tag (updated by deployment)
# - Port 8001 on host (8000 internal) to avoid conflicts with other services
# - Named container for easy management
# - Pulls health check configuration from Containerfile
ExecStart=/usr/bin/podman run \
    --name crimecity3k \
    --rm \
    -p 127.0.0.1:8001:8000 \
    localhost/crimecity3k:production

# Stop container gracefully
ExecStop=/usr/bin/podman stop -t 10 crimecity3k

# Automatically restart on failure
# Container health checks will cause restart if unhealthy
RestartForceExitStatus=0

[Install]
WantedBy=default.target
